# opsweb/supervisor_conf/celery_beat.ini
[program:celery-beat]
command = /ops-data/zp/haha/dev-mysite/bin/python manage.py celery beat
directory = /ops-data/zp/haha/dev-mysite/opsweb_dev
user = root
autostart = true
autorestart = true
stopwaitsecs = 5
redirect_stderr=true
stdout_logfile = /ops-data/zp/haha/dev-mysite/opsweb_dev/logs/supervisor_celery_beat.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10
