# opsweb/conf/supervisor_conf/celery_worker.ini
[program:celery-worker]
command = ***virtualenv***/bin/python manage.py celery worker -l info -n opsweb -B --autoscale=4,2
directory = ***virtualenv***/opsweb
user = root
autostart = true
autorestart = true
stopwaitsecs = 5
redirect_stderr=true
stdout_logfile = ***virtualenv***/opsweb/logs/supervisor_celery_worker.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10

